# systemd/midlab-lis-bridge@.service
[Unit]
Description=MidLab LisBridgeService (instrument %i)
After=network.target mysqld.service
Wants=mysqld.service

[Service]
Type=simple
User=midlab
Group=midlab
WorkingDirectory=/opt/midlab
ExecStart=/usr/bin/python3 -m services.lis_bridge.main --instrument-id %i
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/midlab/lis_bridge_%i.log
StandardError=append:/var/log/midlab/lis_bridge_%i.log

[Install]
WantedBy=multi-user.target
