# systemd/midlab-order-receiver.service
[Unit]
Description=MidLab Order Receiver Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=midlab
Group=midlab
WorkingDirectory=/opt/midlab
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/midlab/.venv/bin/python services/order_receiver/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=midlab-order-receiver

[Install]
WantedBy=multi-user.target
