# systemd/midlab-result-sender.service
[Unit]
Description=MidLab Result Sender Service
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=midlab
Group=midlab
WorkingDirectory=/opt/midlab
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 services/result_sender/main.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=midlab-result-sender

[Install]
WantedBy=multi-user.target
